# woc.service
[Unit]
Description=wifi-offload-concentrator

[Service]
ExecStart=/root/wi-oc/wifi-offload-concentrator.sh

[Install]
WantedBy=default.target
